# msg-service-tools.manifest
<manifest>
	<request>
		<domain name="msg-service"/>
	</request>

	<assign>
		<filesystem path="/etc/rc.d/init.d/msg-server" label="_" exec_label="none"/>
		<filesystem path="/etc/config/sysinfo-message.xml" label="_"/>
	</assign>
</manifest>
