# packaging/msg-server.socket
[Socket]
SocketUser=messaging
SocketGroup=messaging
ListenStream=/tmp/.msgfw_socket
Service=msg-server.service

[Install]
WantedBy=sockets.target
