# packaging/msg-service.service

[Unit]
Description=Start the msg-server

[Service]
ExecStart=/usr/bin/msg-server

[Install]
WantedBy=tizen-middleware.target

